# packages/edgebox/edgebox-0.12.1.tar.gz/edgebox-0.12.1/edgebox/linux_edgebox.service.tmpl
[Unit]
Description=Start up all MobiledgeX Edgeboxes

[Service]
Type=oneshot
ExecStart={{edgebox_comm}} startall
KillMode=process
