# packages/edgebox/edgebox-0.12.1.tar.gz/edgebox-0.12.1/edgebox/linux_edgebox.timer.tmpl
[Unit]
Description=Start up all MobiledgeX Edgeboxes

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
