# packages/edgebox/edgebox-0.12.1.tar.gz/edgebox-0.12.1/edgebox/macos_launchagent.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>{{service_name}}</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<false/>
	<key>StartInterval</key>
	<integer>600</integer>
	<key>ProgramArguments</key>
	<array>
		<string>{{edgebox_comm}}</string>
		<string>startall</string>
	</array>
	<key>StandardOutPath</key>
	<string>{{log_path}}</string>
	<key>StandardErrorPath</key>
	<string>{{log_path}}</string>
</dict>
</plist>
